# systemd/anypc-server.service
# anypc WebSocket inference gateway.
#
# A systemd *user* service installed by `make install` and controlled via the
# `anypc` helper (anypc launch|stop|restart|status|logs). Enabled at boot via
# `systemctl --user enable`; user lingering must be on so the user manager
# starts without a login (`sudo loginctl enable-linger $USER`).

[Unit]
Description=anypc WebSocket inference gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/anypc
Environment=ANYPC_IPC_ENDPOINT=ipc:///tmp/anypc-inference.ipc
ExecStart=/opt/anypc/.venv/bin/anypc-server --host 0.0.0.0 --port 8000

Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
